# assets/systemd/inir.service
[Unit]
Description=iNiR shell
PartOf=graphical-session.target
After=graphical-session.target
Wants=graphical-session.target
# Rate-limit restarts: max 3 attempts in 30 seconds, then stop trying.
# These MUST live in [Unit], not [Service] — systemd ignores them elsewhere.
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Type=simple

# Pre-seed stable environment so the bash wrapper skips detection at boot.
# QT_SCALE_FACTOR must be 1 — shell scaling is handled in QML.
# QT_LOGGING_RULES suppresses noisy Qt/KF warnings.
Environment=QT_SCALE_FACTOR=1
Environment=QT_LOGGING_RULES=quickshell.dbus.properties=false;qt.qml.settings.warning=false;qt.core.qsettings.warning=false;kf.xmlgui=false;kf.coreaddons=false;kf.config.core=false;kf.iconthemes=false

ExecStart=/usr/bin/inir run --session
SuccessExitStatus=143

# process mode: only SIGTERM the main PID so user-spawned apps (terminals,
# media players, etc.) survive shell restarts.  Orphan shell helpers are
# cleaned by ExecStopPost instead.
KillMode=process
KillSignal=SIGTERM

# on-abnormal: restart on signals (SIGSEGV, SIGABRT) but NOT on non-zero
# exit codes.  This respects Quickshell's own crash-loop detector — when
# QS decides it has crashed too many times and exits cleanly, systemd
# won't override that decision by restarting again.
Restart=on-abnormal
RestartSec=5
TimeoutStopSec=15

# Disable coredumps for this service.  Quickshell segfaults on QML hot-reload
# and other internal teardown paths; each dump is 45-90 MB and accumulates fast
# in /var/lib/systemd/coredump.  These are QS-internal crashes, not iNiR bugs —
# debugging them requires a QS developer build, not production coredumps.
LimitCORE=0

# Clean up stale Quickshell runtime entries (by-pid, by-id) left behind
# after crashes.  KillMode=process only kills the main PID, so orphaned
# QS child processes and their runtime artifacts accumulate over time.
ExecStopPost=-/usr/bin/inir cleanup-orphans

# Boost I/O priority so QML/Qt library loading wins over background services
# during cold-boot I/O contention (portals, tray apps, etc.)
IOSchedulingPriority=2

[Install]
WantedBy=graphical-session.target
